# crontabs/ldap
23 5 * * * openldap . /env ; /entrypoint/backup
