# crontabs/ldap
* 5 * * * /entrypoint/backup
* 5 * * * /entrypoint/update_ssl
